# deploy/launchd/com.spacetraders.routing.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.spacetraders.routing</string>
  <key>ProgramArguments</key><array>
    <string>/bin/bash</string><string>/Users/andres.dandrea/projects/spacetraders/gobot/services/routing-service/run.sh</string>
  </array>
  <key>WorkingDirectory</key><string>/Users/andres.dandrea/projects/spacetraders/gobot</string>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/Users/andres.dandrea/projects/spacetraders/gobot/routing.log</string>
  <key>StandardErrorPath</key><string>/Users/andres.dandrea/projects/spacetraders/gobot/routing.log</string>
</dict></plist>
