# deploy/launchd/com.spacetraders.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.spacetraders.daemon</string>
  <key>ProgramArguments</key><array>
    <string>/Users/andres.dandrea/projects/spacetraders/gobot/bin/spacetraders-daemon</string><string>--force</string>
  </array>
  <key>WorkingDirectory</key><string>/Users/andres.dandrea/projects/spacetraders/gobot</string>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>ThrottleInterval</key><integer>10</integer>
  <key>StandardOutPath</key><string>/Users/andres.dandrea/projects/spacetraders/gobot/daemon.log</string>
  <key>StandardErrorPath</key><string>/Users/andres.dandrea/projects/spacetraders/gobot/daemon.log</string>
</dict></plist>
